# templates/centos-7/server/puma.service.erb
[Unit]
<%= render 'server/service/unit.erb' %>

[Service]
Type=forking
<%= render 'server/service/service.erb' %>
PIDFile=<%= @path %>/shared/tmp/pids/puma.pid
Environment=PUMA_CONFIG_FILE=<%= @path %>/shared/config/puma.rb
ExecStart=<%= node['rbenv']['root_path'] %>/shims/bundle exec puma -C ${PUMA_CONFIG_FILE} --daemon
ExecStop=<%= node['rbenv']['root_path'] %>/shims/bundle exec pumactl -F ${PUMA_CONFIG_FILE} stop
ExecReload=<%= node['rbenv']['root_path'] %>/shims/bundle exec pumactl -F ${PUMA_CONFIG_FILE} phased-restart

[Install]
<%= render 'server/service/install.erb' %>
