# templates/centos-7/server/sidekiq.service.erb
[Unit]
<%= render 'server/service/unit.erb' %>

[Service]
Type=simple
<%= render 'server/service/service.erb' %>

PIDFile=<%= @path %>/shared/tmp/pids/sidekiq.pid

Environment=SIDEKIQ_CONFIG_FILE=<%= @path %>/shared/config/sidekiq.yml

ExecStart=/bin/bash -lc 'bundle exec sidekiq -C ${SIDEKIQ_CONFIG_FILE}'

Umask8002

RestartSec=1
Restart=on-failure

[Install]
<%= render 'server/service/install.erb' %>
