# templates/centos-7/server/sidekiq.service.erb
[Unit]
<%= render 'server/service/unit.erb' %>

[Service]
Type=forking
<%= render 'server/service/service.erb' %>
Restart=on-abort
PIDFile=<%= @path %>/shared/tmp/pids/sidekiq.pid
Environment=SIDEKIQ_CONFIG_FILE=<%= @path %>/shared/config/sidekiq.yml
ExecStart=<%= node['rbenv']['root_path'] %>/shims/bundle exec sidekiq -C ${SIDEKIQ_CONFIG_FILE} --daemon
ExecStop=<%= node['rbenv']['root_path'] %>/shims/bundle exec sidekiqctl stop ${<%= @path %>/shared/tmp/pids/sidekiq.pid}

[Install]
<%= render 'server/service/install.erb' %>
