# templates/centos-7/server/sidekiq.service.erb
[Unit]
<%= render 'server/service/unit.erb' %>

[Service]
Type=simple
<%= render 'server/service/service.erb' %>
PIDFile=<%= @opts[:path] %>/shared/tmp/pids/sidekiq.pid
Environment=SIDEKIQ_CONFIG_FILE=<%= @opts[:path] %>/shared/config/sidekiq.yml
ExecStart=<%= node['rbenv']['root_path'] %>/shims/bundle exec sidekiq -C ${SIDEKIQ_CONFIG_FILE}
ExecStop=<%= node['rbenv']['root_path'] %>/shims/bundle exec sidekiqctl stop ${<%= @opts[:path] %>/shared/tmp/pids/sidekiq.pid}

[Install]
<%= render 'server/service/install.erb' %>
