# templates/default/etc/profile.d/temp-folder.sh.erb
export TMPDIR="<%= node['rails']['azure']['mnt'] %>"
export TMP="<%= node['rails']['azure']['mnt'] %>"
