# templates/ubuntu-16/ethminer.service.erb
# Maintained by Chef.  Local changes will be overwritten

[Unit]
Description=Ethminer client
Requires=network.target

[Service]
Type=forking
User=%I
ExecStart=ethminer -G -S <%= @pool_address %> -u <%= @wallet_address %>
Restart=on-failure

[Install]
WantedBy=multi-user.target
